# provision/supervisor/files/deployment.conf
[program:deployment]
command=/bin/bash /home/deploy/deployment.sh
user=deploy
redirect_stderr=true
autostart=false
directory=/home/deploy