# provision/supervisor/files/logviewer.conf
[program:logviewer]
command=python3 /home/deploy/logviewer/app.py
user=deploy
redirect_stderr=true
autostart=true
directory=/home/deploy